# overlay/system/etc/profile.d/powervr-fix.sh
#!/bin/bash

# This fixes the black lines experienced on PowerVR devices.
# Thanks to Alfred and NotKit for finding this.

export UC_SHAPE_MIPMAPS=1
