# overlay/system/usr/lib/systemd/user/pulseaudio.service.d/ubuntu-touch.conf
[Service]
# Use the Touch specific Pulseaudio script file
#Environment=PULSE_SCRIPT=/etc/pulse/touch.pa

# Force server-side corking when stalled (playback)
Environment=PULSE_PLAYBACK_CORK_STALLED=1
